# deploy/systemd/whisper-asr-update.timer
# /etc/systemd/system/whisper-asr-update.timer
# Triggers whisper-asr-update.service periodically.
#
# To adjust the check interval after installation:
#   sudo systemctl edit whisper-asr-update.timer
#   # Add:
#   #   [Timer]
#   #   OnUnitActiveSec=600
#   sudo systemctl daemon-reload
#   sudo systemctl restart whisper-asr-update.timer

[Unit]
Description=Whisper ASR Auto-Update Timer
Documentation=https://github.com/revoice-resonance/asr
Requires=whisper-asr-update.service

[Timer]
# Check for updates every 5 minutes, with 30s randomized delay to avoid
# thundering herd against GitHub
OnUnitActiveSec=300
RandomizedDelaySec=30

# If a check was missed (e.g., server was down), run immediately on boot
Persistent=true

[Install]
WantedBy=timers.target
