# deploy/systemd/whisper-asr.service
# /etc/systemd/system/whisper-asr.service
# Whisper ASR API — main server service
#
# Installation (via deploy.sh):
#   bash deploy.sh install-systemd
#
# Manual installation:
#   sudo cp deploy/systemd/whisper-asr.service /etc/systemd/system/
#   sudo cp deploy/systemd/whisper-asr-update.service /etc/systemd/system/
#   sudo cp deploy/systemd/whisper-asr-update.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now whisper-asr.service
#   sudo systemctl enable --now whisper-asr-update.timer

[Unit]
Description=Whisper ASR API Server
Documentation=https://github.com/revoice-resonance/asr
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=whisper
Group=whisper
WorkingDirectory=/opt/whisper-asr

# Load environment from .env file (systemd v240+ supports # comments)
# The leading "-" means "ignore if file missing"
EnvironmentFile=-/opt/whisper-asr/.env

# Start via uvicorn in the venv
# ${VAR:-default} syntax is systemd's own variable expansion from EnvironmentFile
ExecStart=/opt/whisper-asr/venv/bin/python -m uvicorn app.main:app \
    --host ${HOST:-0.0.0.0} \
    --port ${PORT:-8080} \
    --log-level ${LOG_LEVEL:-info} \
    --workers 1

# Graceful shutdown: SIGTERM with 35s timeout before SIGKILL
TimeoutStopSec=35
KillMode=mixed
KillSignal=SIGTERM

# Restart policy — restart on crash, not on clean exit
Restart=on-failure
RestartSec=10

# GPU needs memlock for CUDA
LimitMEMLOCK=infinity
LimitNOFILE=65536

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whisper-asr

# Sandbox hardening
NoNewPrivileges=yes
# NOTE: PrivateTmp is intentionally disabled — TEMP_DIR (from .env) is used
# for audio temp files and must be on the shared filesystem.
ProtectSystem=strict
ProtectHome=yes
# __READ_WRITE_PATHS__ is replaced by deploy.sh install-systemd at install time
ReadWritePaths=__READ_WRITE_PATHS__
ReadOnlyPaths=__READ_ONLY_PATHS__
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
